# docker/s6/etc/services.d/app/run
#!/usr/bin/with-contenv sh

python3.8 ./proxy_switcher.py

